# deb/postinst
#!/bin/sh
# postinst script for Comdb2
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        adduser --quiet --system --group --shell /bin/bash comdb2
        for dir in /opt/bb/var /opt/bb/var/log ; do
            [ ! -d $dir ] && mkdir -p $dir 
            chmod 755 $dir
        done
        for dir in /opt/bb/var/run; do
            [ ! -d $dir ] && mkdir -p $dir 
            chmod 777 $dir
        done
        for dir in /opt/bb/var/cdb2/ /opt/bb/etc/cdb2 /opt/bb/var/log/cdb2 /opt/bb/var/lib/cdb2 /opt/bb/etc/cdb2/config/comdb2.d/ /opt/bb/var/log/cdb2_supervisor /opt/bb/etc/cdb2_supervisor/conf.d/; do
            [ ! -d $dir ] && mkdir -p $dir 
            chown comdb2:comdb2 $dir
            chmod 775 $dir
        done
        echo 'PATH=$PATH:/opt/bb/bin' >> /home/comdb2/.profile
        chmod +x /home/comdb2/.profile
        chown comdb2:comdb2 /home/comdb2/.profile
        if [ ! -e /.dockerenv ]; then
            systemctl daemon-reload
            systemctl stop pmux
            systemctl start pmux
            systemctl stop cdb2sockpool
            systemctl start cdb2sockpool
            set +e
            dpkg -l | grep -q "ii supervisor"
            rc=$?
            set -e
            if [ $rc -eq 0 ]; then
                systemctl enable supervisor_cdb2
                systemctl start supervisor_cdb2
            fi
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
